# device-redis-el9/src/main/resources/etc/systemd/system/redis.service.d/local.conf
[Unit]
PartOf=redis-autodiscovery.service

